# templates/home-security-system.service.template
[Unit]\\nDescription=Home security system\\nAfter=network.target\\n[Service]\\nType=simple\\nRemainAfterExit=yes\\nRestart=always\\nExecStart=python3 $(pwd)/app.py\\nWorkingDirectory=$(pwd)\\nUser=$(whoami)\\n[Install]\\nWantedBy=multi-user.target\\n
